# deploy/pineapple-graph-rebuild.timer
[Unit]
Description=Pineapple OS — nightly vault graph rebuild timer

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
